# deploy/supervisord.conf
[unix_http_server]
file=/opt/corporate-rag-mvp/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/opt/corporate-rag-mvp/var/log/supervisord.log
pidfile=/opt/corporate-rag-mvp/var/run/supervisord.pid
childlogdir=/opt/corporate-rag-mvp/var/log
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///opt/corporate-rag-mvp/var/run/supervisor.sock

[program:embedder]
command=/opt/corporate-rag-mvp/.venv/bin/python /opt/corporate-rag-mvp/scripts/run_embedder.py
directory=/opt/corporate-rag-mvp
user=ragapp
autostart=true
autorestart=true
startsecs=15
stopasgroup=true
killasgroup=true
stdout_logfile=/opt/corporate-rag-mvp/var/log/embedder.stdout.log
stderr_logfile=/opt/corporate-rag-mvp/var/log/embedder.stderr.log
environment=HOME="/home/ragapp",PYTHONPATH="/opt/corporate-rag-mvp/src",HF_HOME="/opt/corporate-rag-mvp/var/cache/huggingface",TRANSFORMERS_CACHE="/opt/corporate-rag-mvp/var/cache/huggingface"

[program:api]
command=/opt/corporate-rag-mvp/.venv/bin/python /opt/corporate-rag-mvp/scripts/run_api.py
directory=/opt/corporate-rag-mvp
user=ragapp
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/opt/corporate-rag-mvp/var/log/api.stdout.log
stderr_logfile=/opt/corporate-rag-mvp/var/log/api.stderr.log
environment=HOME="/home/ragapp",PYTHONPATH="/opt/corporate-rag-mvp/src"

[program:telegram-bot]
command=/opt/corporate-rag-mvp/.venv/bin/python /opt/corporate-rag-mvp/scripts/run_bot.py
directory=/opt/corporate-rag-mvp
user=ragapp
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/opt/corporate-rag-mvp/var/log/bot.stdout.log
stderr_logfile=/opt/corporate-rag-mvp/var/log/bot.stderr.log
environment=HOME="/home/ragapp",PYTHONPATH="/opt/corporate-rag-mvp/src"
